# PiNetOS/services/minima.service
[Unit]
Description=Minima Blockchain Node
After=network.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/java -jar /opt/minima/minima.jar
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
